# examples/com.jellygenie.chrome.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jellygenie.chrome</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Google Chrome.app/Contents/MacOS/Google Chrome</string>
        <string>--user-data-dir=/Users/YOURNAME/.jellygenie/browser-profile</string>
        <string>--remote-debugging-port=9222</string>
        <string>--remote-debugging-address=127.0.0.1</string>
        <string>--no-first-run</string>
        <string>--no-default-browser-check</string>
        <string>--restore-last-session</string>
        <string>--disable-features=ChromeWhatsNewUI</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/YOURNAME</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/YOURNAME</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>15</integer>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/tmp/jellygenie-logs/chrome.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/jellygenie-logs/chrome.err.log</string>
</dict>
</plist>
